# conf/supervisor.conf
[program:captcha]
command=/captcha/env/bin/flask run
directory=/captcha/src
user=root
autostart=true
autorestart=true
redirect_stderr=True
environment=FLASK_APP=server.py
